# packaging/systemd/vps-guard-control.service
[Unit]
Description=VPSGuard control plane and loopback operations console
After=network.target

[Service]
Type=simple
User=vps-guard
Group=vps-guard
ExecStart=/usr/local/bin/vps-guard-control
Environment=VPS_GUARD_CONFIG=/etc/vps-guard/config.toml
Environment=VPS_GUARD_STATE=/var/lib/vps-guard/state.json
Environment=RUST_LOG=vps_guard_control=info,guard_control=info,guard_agent=info,guard_provider=info,guard_system=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vps-guard-control
LogRateLimitIntervalSec=30s
LogRateLimitBurst=2000
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
Restart=on-failure
RestartSec=3s
RuntimeDirectory=vps-guard
RuntimeDirectoryMode=0750
StateDirectory=vps-guard
StateDirectoryMode=0750
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/run/vps-guard /var/lib/vps-guard
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
UMask=0027
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
SystemCallArchitectures=native
TasksMax=256

[Install]
WantedBy=multi-user.target
